# docker/choralstorm/supervisor/kafka-server.conf
[program:KafkaServer]
command=java -cp /kafkaserver/kafkaserver.jar KafkaServer
autostart=true
autorestart=true
