# docker/choralstorm/supervisor/storm-supervisor.conf
[program:storm-supervisor]
command=/usr/bin/start-storm-supervisor.sh
autostart=true
autorestart=true
