# docker/choralstorm/supervisor/storm-ui.conf
[program:storm-ui]
command=/usr/bin/start-storm-ui.sh
autostart=true
autorestart=true
