# docker/choralstorm/supervisor/storm-nimbus.conf
[program:storm-nimbus]
command=/usr/bin/start-storm-nimbus.sh
autostart=true
autorestart=true
